# install/cron/app.cron
MAILTO=""

# Application logs
0 */6 * * * root php /duktig.microservice.1/cli/exec.php generateLogStats
5 */6 * * * root php /duktig.microservice.1/cli/exec.php archiveLogFiles

# Databases Backup
0 0 * * * root php /duktig.microservice.1/cli/exec.php db-backup
